# deploy/systemd/levant-gateway.service
# Levant AI Node gateway / BFF (Fastify + WebSocket).
# Runs as your existing login user. If your user is not "opt", change User/Group
# below and the .nvm path in ExecStart to match your home directory.

[Unit]
Description=Levant AI gateway (Node/Fastify)
After=network.target levant-ai.service
Wants=levant-ai.service

[Service]
Type=simple
User=opt
Group=opt
WorkingDirectory=/opt/levant/services/gateway
EnvironmentFile=/opt/levant/services/gateway/.env
# Login shell sources nvm and runs whatever Node version you installed —
# version-proof (no hardcoded vXX path). Needs ProtectHome=read-only (below)
# so the nvm runtime under /home is readable.
ExecStart=/bin/bash -lc 'source /home/opt/.nvm/nvm.sh && exec node dist/index.js'
Restart=on-failure
RestartSec=3

# --- Isolation / resource guardrails ---
MemoryMax=256M
CPUQuota=40%
TasksMax=128
Nice=10
IOWeight=50

# --- Hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# read-only (not true) so the nvm Node runtime under /home/<user> is readable.
ProtectHome=read-only
ReadWritePaths=/opt/levant/services/gateway
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
